# cachyos/updates.service
[Unit]
Description=System updates

[Service]
Type=oneshot
ExecStart=%h/bin/updates.sh
