# cachyos/updates.timer
[Unit]
Description=System updates

[Timer]
OnCalendar=*-*-* 00:30:00
Persistent=true

[Install]
WantedBy=timers.target
